# rcc/dashboard/squirrelchat-server/squirrelchat-server/squirrelchat-server.service
[Unit]
Description=SquirrelChat Server v2 (Rust/Axum)
After=network.target

[Service]
Type=simple
User=jkh
WorkingDirectory=/home/jkh/.openclaw/workspace/squirrelchat
ExecStart=/home/jkh/rockyandfriends/rcc/dashboard/target/release/squirrelchat-server
Environment=SQUIRRELCHAT_PORT=8793
Environment=SQUIRRELCHAT_DB=/home/jkh/.openclaw/workspace/squirrelchat/squirrelchat.db
Environment=RUST_LOG=info
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
