# rcc/crush-server/crush-server.service
[Unit]
Description=Crush Server — HTTP/SSE bridge for charmbracelet/crush (with claude-code failover)
After=network.target

[Service]
Type=simple
User=jkh
WorkingDirectory=/home/jkh/rockyandfriends/rcc/crush-server
ExecStart=/usr/bin/env node /home/jkh/rockyandfriends/rcc/crush-server/index.mjs
Restart=on-failure
RestartSec=5
Environment=NODE_ENV=production
Environment=CRUSH_SERVER_PORT=8793
Environment=CRUSH_BIN=/home/jkh/.npm-global/bin/crush
# claude-code binary (failover source); set to empty to disable claude-code attempts
Environment=CLAUDE_BIN=claude
# RCC URL + token for bus notifications on throttle events
Environment=RCC_URL=http://localhost:8789
# CRUSH_ONLY=1  → always use crush, skip claude-code
# CLAUDE_THROTTLED=1 → force crush without trying claude-code this boot
# EnvironmentFile=/home/jkh/.config/crush/env  (for ANTHROPIC_API_KEY etc.)

[Install]
WantedBy=multi-user.target
