# deploy/systemd/acc-server.service
[Unit]
Description=ACC API Server (Rust/Axum)
Documentation=https://github.com/{OWNER}/CCC
After=network.target

[Service]
Type=simple
User=AGENT_USER
WorkingDirectory=AGENT_HOME/Src/CCC

ExecStart=/usr/local/bin/acc-server

# Core config
Environment=ACC_PORT=8789
Environment=ACC_DATA_DIR=AGENT_HOME/.acc/data
Environment=RUST_LOG=acc_server=info

# Load per-node secrets
EnvironmentFile=-AGENT_HOME/.acc/.env

# Restart on crash, not clean exits
Restart=on-failure
RestartSec=5s

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=acc-server

[Install]
WantedBy=multi-user.target
