# deploy/systemd/acc-agent.service
[Unit]
Description=ACC Agent Supervisor
Documentation=https://github.com/{OWNER}/ACC
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=AGENT_USER
EnvironmentFile=-AGENT_HOME/.acc/.env
Environment=PATH=AGENT_HOME/.acc/bin:AGENT_HOME/.local/bin:/usr/local/bin:/usr/bin:/bin
ExecStart=AGENT_HOME/.acc/bin/acc-agent supervise
Restart=on-failure
RestartSec=5s
StandardOutput=journal
StandardError=journal
SyslogIdentifier=acc-agent

[Install]
WantedBy=multi-user.target
