# deploy/systemd/acc-hermes-worker.service
[Unit]
Description=ACC Hermes Worker (polls queue for hermes-tagged tasks)
Documentation=https://github.com/{OWNER}/ACC
After=network-online.target acc-bus-listener.service
Wants=network-online.target
# Only enable on nodes where `hermes` binary is present

[Service]
Type=simple
User=AGENT_USER
EnvironmentFile=-AGENT_HOME/.acc/.env
Environment=PATH=AGENT_HOME/.local/bin:AGENT_HOME/.venv/bin:/usr/local/bin:/usr/bin:/bin
ExecStart=AGENT_HOME/.acc/bin/acc-agent hermes --poll
Restart=on-failure
RestartSec=30s

StandardOutput=journal
StandardError=journal
SyslogIdentifier=acc-hermes-worker

[Install]
WantedBy=multi-user.target
