# deploy/systemd/acc-github-sync.timer
[Unit]
Description=ACC GitHub issues to beads/fleet sync timer
Requires=acc-github-sync.service

[Timer]
OnBootSec=2min
OnUnitActiveSec=5min
RandomizedDelaySec=30
Persistent=true

[Install]
WantedBy=timers.target
