# deploy/systemd/acc-queue-worker.service
[Unit]
Description=ACC Queue Worker (autonomous task executor)
Documentation=https://github.com/jordanhubbard/CCC
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=AGENT_USER
EnvironmentFile=-AGENT_HOME/.acc/.env
Environment=PATH=AGENT_HOME/.local/bin:/usr/local/bin:/usr/bin:/bin
ExecStart=/usr/bin/python3 AGENT_HOME/.acc/workspace/deploy/queue-worker.py
Restart=always
RestartSec=10s
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
