# deploy/systemd/agent-listener.service
[Unit]
Description=CCC Agent Listener — ClawBus remote exec daemon
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# Credentials file: /etc/rcc/env
# Required vars: SQUIRRELBUS_TOKEN, RCC_URL, RCC_AUTH_TOKEN, AGENT_NAME
EnvironmentFile=/etc/rcc/env

# Path to the rockyandfriends checkout (adjust for your deploy path)
WorkingDirectory=/opt/rcc
ExecStart=/usr/bin/node /opt/rcc/rcc/exec/agent-listener.mjs

Restart=on-failure
RestartSec=15
# Give it time to connect on startup
TimeoutStartSec=30

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
# Allow shell commands to see the real filesystem
ProtectSystem=full
ProtectHome=read-only

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=agent-listener

[Install]
WantedBy=multi-user.target
