# deploy/systemd/ccc-server.service
[Unit]
Description=CCC API Server (Rust/Axum)
Documentation=https://github.com/jordanhubbard/rockyandfriends
After=network.target

[Service]
Type=simple
User=jkh
WorkingDirectory=/home/jkh/Src/CCC

# Binary built from ccc/dashboard/ccc-server via:
#   cd ccc/dashboard && cargo build -p ccc-server --release
ExecStart=/usr/local/bin/ccc-server

# Core config
Environment=CCC_PORT=8789
Environment=CCC_DATA_DIR=/home/jkh/.ccc/data
Environment=RUST_LOG=ccc_server=info

# Restart on crash, not clean exits
Restart=on-failure
RestartSec=5s

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ccc-server

[Install]
WantedBy=multi-user.target
