# deploy/systemd/ccc-dashboard.service
[Unit]
Description=CCC Dashboard Server (nginx + WASM SPA)
Documentation=https://github.com/jordanhubbard/rockyandfriends
After=network.target ccc-server.service
Wants=ccc-server.service

[Service]
Type=simple
User=jkh
WorkingDirectory=/home/jkh/Src/CCC

# Serves the Leptos WASM SPA from ccc/dashboard/dist/
# Build with: cd ccc/dashboard && make release-ui
ExecStart=/home/jkh/Src/CCC/ccc/dashboard/target/release/dashboard-server

Environment=RCC_DASHBOARD_PORT=8788
Environment=RCC_URL=http://127.0.0.1:8789
Environment=DASHBOARD_DIST=/home/jkh/Src/CCC/ccc/dashboard/dist
Environment=CCC_AGENT_TOKEN=

Restart=on-failure
RestartSec=5s

StandardOutput=journal
StandardError=journal
SyslogIdentifier=ccc-dashboard

[Install]
WantedBy=multi-user.target
