# deploy/systemd/ccc-agent.service
[Unit]
Description=CCC Agent Pull Service
After=network.target
Wants=network-online.target

[Service]
Type=oneshot
User=AGENT_USER
EnvironmentFile=-AGENT_HOME/.acc/.env
ExecStart=/bin/bash AGENT_HOME/.acc/workspace/deploy/agent-pull.sh
StandardOutput=append:AGENT_HOME/.acc/logs/pull.log
StandardError=append:AGENT_HOME/.acc/logs/pull.log

[Install]
WantedBy=multi-user.target
