# deploy/systemd/clawfs-sparky.service
[Unit]
Description=ClawFS (JuiceFS) mount at ~/clawfs
After=network-online.target consul.service
Wants=network-online.target

[Service]
Type=forking
# Redis address via Consul DNS — requires migration 0010 (*.consul DNS)
ExecStart=AGENT_HOME/.local/bin/juicefs mount --background --cache-dir /tmp/jfscache redis://ccc-server.service.consul:6379/1 AGENT_HOME/clawfs
ExecStop=AGENT_HOME/.local/bin/juicefs umount AGENT_HOME/clawfs
Restart=on-failure
RestartSec=30

[Install]
WantedBy=default.target
