# deploy/systemd/clawbus-natasha.service
[Unit]
Description=ClawBus receive sidecar for Natasha
After=network.target

[Service]
Type=simple
User=jkh
ExecStart=/usr/bin/node /home/jkh/.openclaw/workspace/clawbus/receive-server.mjs
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=clawbus-natasha

[Install]
WantedBy=multi-user.target
