# deploy/systemd/consul.service
[Unit]
Description=Consul Service Discovery Agent (CCC)
Documentation=https://www.consul.io/
# Start after network and Tailscale — Consul binds to Tailscale IP
After=network-online.target tailscaled.service
Wants=network-online.target tailscaled.service

[Service]
Type=notify
User=AGENT_USER
ExecStart=AGENT_HOME/.ccc/consul/bin/consul agent -config-file=AGENT_HOME/.ccc/consul/consul.hcl
ExecReload=/bin/kill -HUP $MAINPID
KillMode=process
KillSignal=SIGTERM
Restart=on-failure
RestartSec=5s
LimitNOFILE=65536

StandardOutput=journal
StandardError=journal
SyslogIdentifier=consul

[Install]
WantedBy=multi-user.target
