# deploy/systemd/mm-bridge.service
[Unit]
Description=Mattermost <-> SquirrelChat bridge
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/rockyandfriends/services/mm-bridge
ExecStart=/usr/bin/node bridge.mjs
Restart=on-failure
RestartSec=10
EnvironmentFile=/etc/mm-bridge.env
StandardOutput=journal
StandardError=journal
SyslogIdentifier=mm-bridge

[Install]
WantedBy=multi-user.target
