# deploy/launchd/com.acc.hermes-worker.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.acc.hermes-worker</string>

    <key>ProgramArguments</key>
    <array>
        <string>AGENT_HOME/.acc/bin/acc-agent</string>
        <string>hermes</string>
        <string>--poll</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>AGENT_HOME/.local/bin:AGENT_HOME/.venv/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>

    <!-- 30s throttle: hermes tasks are long-running, don't restart-storm -->
    <key>ThrottleInterval</key>
    <integer>30</integer>

    <key>StandardOutPath</key>
    <string>AGENT_HOME/.acc/logs/hermes-worker.log</string>
    <key>StandardErrorPath</key>
    <string>AGENT_HOME/.acc/logs/hermes-worker.log</string>
</dict>
</plist>
