# deploy/launchd/com.acc.nvidia-proxy.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.acc.nvidia-proxy</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-c</string>
        <string>exec "$(PATH=/opt/homebrew/bin:/usr/local/bin:$PATH which python3)" AGENT_HOME/.acc/workspace/deploy/nvidia-proxy.py --port 9099</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>AGENT_HOME/.acc/logs/nvidia-proxy.log</string>
    <key>StandardErrorPath</key>
    <string>AGENT_HOME/.acc/logs/nvidia-proxy.log</string>
    <key>ThrottleInterval</key>
    <integer>5</integer>
</dict>
</plist>
