# deploy/launchd/com.acc.consul.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.acc.consul</string>

    <key>ProgramArguments</key>
    <array>
        <string>AGENT_HOME/.acc/consul/bin/consul</string>
        <string>agent</string>
        <string>-config-file=AGENT_HOME/.acc/consul/consul.hcl</string>
    </array>

    <!-- Run at login and keep alive on crash -->
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>AGENT_HOME/.acc/logs/consul.log</string>
    <key>StandardErrorPath</key>
    <string>AGENT_HOME/.acc/logs/consul.log</string>

    <!-- Throttle respawn on persistent failure -->
    <key>ThrottleInterval</key>
    <integer>5</integer>
</dict>
</plist>
