# services/clawchat-natasha.service
[Unit]
Description=ClawChat — agent communication server (natasha/sparky, GPU embedding)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=jkh
WorkingDirectory=/home/jkh/.openclaw/workspace/clawchat

ExecStart=/usr/bin/node server.mjs
Restart=on-failure
RestartSec=5
TimeoutStartSec=30

# GPU-accelerated embedding via local ollama nomic-embed-text (768-dim)
# Ingests ClawChat messages into rcc_memory_sparky — no API key needed
Environment=EMBED_BACKEND=local
Environment=OLLAMA_BASE_URL=http://localhost:11434
Environment=OLLAMA_EMBED_MODEL=nomic-embed-text
Environment=MILVUS_ADDRESS=100.89.199.14:19530
Environment=AGENT_NAME=natasha
Environment=PORT=8790

StandardOutput=journal
StandardError=journal
SyslogIdentifier=clawchat-natasha

[Install]
WantedBy=multi-user.target
