# services/openclaw-container.conf
; OpenClaw Gateway supervisord config for container vLLM agents
; Place at /etc/supervisor/conf.d/openclaw.conf or ~/supervisord.conf.d/openclaw.conf
; Install: supervisorctl reread && supervisorctl update

[program:openclaw-gateway]
; Poll-only mode — no inbound port required; outbound SSE connection to RCC
command=/usr/bin/openclaw gateway start --poll-only --poll-interval 30 --heartbeat-interval 60
directory=/root
autostart=true
autorestart=true
startretries=5
startsecs=10
redirect_stderr=true
stdout_logfile=/var/log/openclaw-gateway.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=3
; Pass through env vars from venv / .env file
environment=HOME="/root",OPENCLAW_GATEWAY_MODE="poll",RCC_AGENT_TOKEN="%(ENV_RCC_AGENT_TOKEN)s",RCC_URL="%(ENV_RCC_URL)s"
