# services/rcc-dashboard.service
[Unit]
Description=RCC Dashboard Server (Rust/WASM)
Documentation=https://github.com/jordanhubbard/rockyandfriends
After=network.target

[Service]
Type=simple
User=jkh
WorkingDirectory=/home/jkh/Src/rockyandfriends/rcc/dashboard

Environment=RCC_DASHBOARD_PORT=8788
Environment=DASHBOARD_DIST=/home/jkh/Src/rockyandfriends/rcc/dashboard/dist
Environment=RCC_URL=https://api.yourmom.photos
Environment=OPERATOR=jkh

ExecStartPre=/bin/sh -c 'test -f /home/jkh/Src/rockyandfriends/rcc/dashboard/target/release/dashboard-server || (echo "dashboard-server binary not found — run: cd ~/Src/rockyandfriends/rcc/dashboard && cargo build -p dashboard-server --release" && exit 1)'
ExecStart=/home/jkh/Src/rockyandfriends/rcc/dashboard/target/release/dashboard-server

Restart=on-failure
RestartSec=5
StartLimitBurst=5
StartLimitIntervalSec=60

StandardOutput=journal
StandardError=journal
SyslogIdentifier=rcc-dashboard

[Install]
WantedBy=multi-user.target
