# deploy/supervisor/acc-queue-worker.conf
[program:acc-queue-worker]
command=/usr/bin/python3 AGENT_HOME/.acc/workspace/deploy/queue-worker.py
directory=AGENT_HOME
user=AGENT_USER
autostart=false
autorestart=true
startsecs=10
startretries=5
redirect_stderr=true
stdout_logfile=AGENT_HOME/.acc/logs/queue-worker.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
environment=HOME="AGENT_HOME"
