# deploy/acc-repo-sync.timer
[Unit]
Description=ACC repo sync timer — runs every 30 minutes
Requires=acc-repo-sync.service

[Timer]
OnBootSec=5min
OnUnitActiveSec=30min
RandomizedDelaySec=60
Persistent=true

[Install]
WantedBy=timers.target
