# deploy/ccc-repo-sync.service
[Unit]
Description=CCC repo sync — pull/commit/push shared ClawFS repo
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=jkh
Group=jkh
WorkingDirectory=/home/jkh
ExecStart=/bin/bash %h/.ccc/workspace/deploy/ccc-repo-sync.sh
EnvironmentFile=-%h/.ccc/.env
StandardOutput=append:%h/.ccc/logs/repo-sync.log
StandardError=append:%h/.ccc/logs/repo-sync.log

# Safety: don't let a stuck git command run forever
TimeoutStartSec=120
