# ccc/dashboard/ccc-dashboard.service
[Unit]
Description=RCC WASM Dashboard Server (Rust/Axum)
Documentation=https://github.com/jordanhubbard/rockyandfriends
After=network.target

[Service]
Type=simple
User=jkh
WorkingDirectory=/home/jkh/Src/rockyandfriends/rcc/dashboard

# Binary built from rcc/dashboard/dashboard-server
ExecStart=/home/jkh/Src/rockyandfriends/rcc/dashboard/target/release/dashboard-server

# Environment
Environment=RCC_DASHBOARD_PORT=8788
Environment=RCC_URL=http://146.190.134.110:8789
Environment=SC_URL=http://localhost:8790
Environment=DASHBOARD_DIST=/home/jkh/Src/rockyandfriends/rcc/dashboard/dist
Environment=RCC_AGENT_TOKEN=wq-5dcad756f6d3e345c00b5cb3dfcbdedb
Environment=OPERATOR_HANDLE=jkh

# Restart on crash (but not clean exits)
Restart=on-failure
RestartSec=5s

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=rcc-dashboard

[Install]
WantedBy=multi-user.target
