# squirrelbus/squirrelbus-natasha.service
[Unit]
Description=SquirrelBus receive sidecar for Natasha
After=network.target

[Service]
Type=simple
User=jkh
ExecStart=/usr/bin/node /home/jkh/.openclaw/workspace/squirrelbus/receive-server.mjs
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=squirrelbus-natasha

[Install]
WantedBy=multi-user.target
