# rcc/clawfs-sync/clawfs-sync.service
[Unit]
Description=AgentFS Sync — MinIO workspace sync daemon
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=AGENT_USER
EnvironmentFile=-/home/AGENT_USER/.rcc/.env
ExecStart=/usr/local/bin/clawfs-sync
Restart=always
RestartSec=15
StandardOutput=append:/home/AGENT_USER/.rcc/logs/clawfs-sync.log
StandardError=append:/home/AGENT_USER/.rcc/logs/clawfs-sync.log

[Install]
WantedBy=multi-user.target
